# deploy/launchd/LaunchDaemons/com.pqw.mp4dl-app.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.pqw.mp4dl-app</string>
    <key>UserName</key>
    <string>jackyhu</string>
    <key>GroupName</key>
    <string>staff</string>
    <key>ProgramArguments</key>
    <array>
      <string>/Volumes/STORAGE_Jackyhu/code/x_video/.venv/bin/python</string>
      <string>server.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Volumes/STORAGE_Jackyhu/code/x_video</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Volumes/STORAGE_Jackyhu/code/x_video/logs/mp4dl-app.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Volumes/STORAGE_Jackyhu/code/x_video/logs/mp4dl-app.err.log</string>
  </dict>
</plist>
